# deploy/yt-api.service
# yours-truly API (Fastify). Install as /etc/systemd/system/yt-api.service.
# Runs from the deploy clone that yt-deploy.sh builds; restarted on each deploy.
[Unit]
Description=yours-truly API (Fastify)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ytapp
Group=ytapp
WorkingDirectory=/opt/yours-truly/apps/api
EnvironmentFile=/etc/yours-truly/api.env
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=2
Nice=10
# Containment: a web-tier RCE stays an unprivileged, sandboxed process.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full

[Install]
WantedBy=multi-user.target
