# deploy/yt-web.service
# yours-truly web (Next.js SSR). Install as /etc/systemd/system/yt-web.service.
# Runs `next start` from the deploy clone; restarted on each deploy.
[Unit]
Description=yours-truly web (Next.js SSR)
After=network-online.target yt-api.service
Wants=network-online.target

[Service]
Type=simple
User=ytapp
Group=ytapp
WorkingDirectory=/opt/yours-truly/apps/web
Environment=NODE_ENV=production
# Do NOT pass `-H 127.0.0.1`. An explicit bind host + nginx's `X-Forwarded-Proto:
# https` makes Next build ABSOLUTE redirects/rewrites from the bind host (e.g.
# https://localhost:3000/en) instead of the public Host, which 500s / dead-redirects
# every locale path behind the proxy (next-intl). Binding all interfaces is safe:
# the SG blocks :3000 from the internet, so nginx stays the only ingress. We trade
# the -H defense-in-depth (now the SG alone) for working i18n redirects.
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3000
Restart=on-failure
RestartSec=2
Nice=10
# Containment: a web-tier RCE stays an unprivileged, sandboxed process.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full

[Install]
WantedBy=multi-user.target
