# deploy/yt-web.service
# yours-truly web (Next.js SSR). Install as /etc/systemd/system/yt-web.service.
# Runs `next start` from the deploy clone; restarted on each deploy.
[Unit]
Description=yours-truly web (Next.js SSR)
After=network-online.target yt-api.service
Wants=network-online.target

[Service]
Type=simple
User=ytapp
Group=ytapp
WorkingDirectory=/opt/yours-truly/apps/web
Environment=NODE_ENV=production
# -H 127.0.0.1: only nginx on this host may reach the app server directly —
# the SG already filters :3000, this removes the dependence on it.
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3000 -H 127.0.0.1
Restart=on-failure
RestartSec=2
Nice=10
# Containment: a web-tier RCE stays an unprivileged, sandboxed process.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full

[Install]
WantedBy=multi-user.target
